# systemd/jurt.service
[Unit]

[Service]
ExecStart=sh /root/jurt.sh

[Install]
WantedBy=default.target
